# systemd/homelab-db-provisioner-oneshot.service
[Unit]
Description=Homelab DB Provisioner (one-shot)
Documentation=https://github.com/Harnish/homelab_db_provisioner
After=network.target

[Service]
Type=oneshot
User=db-provisioner
Group=db-provisioner
EnvironmentFile=/etc/homelab-db-provisioner/env
ExecStart=/usr/local/bin/homelab-db-provisioner
StandardOutput=journal
StandardError=journal
SyslogIdentifier=homelab-db-provisioner

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/etc/homelab-db-provisioner /var/lib/homelab-db-provisioner
PrivateTmp=true

[Install]
WantedBy=multi-user.target
